# conf/indeliblevsan.conf
description "Indelible Web Access"

start on (started indeliblefs and runlevel [2345] )

stop on runlevel [016]

pre-start script
	logger "Indelible Virtual SAN start"
end script

script
	exec /usr/bin/launchIndelibleVirtualSAN | logger 2>&1
end script
